# config/systemd/conecta-warmup.service
[Unit]
Description=Conecta Plus Cache Warmup Service
Documentation=https://github.com/conecta-plus/docs
After=docker.service
Requires=docker.service
ConditionPathExists=/opt/conecta-plus/scripts/warmup.py

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/conecta-plus
Environment=BACKEND_URL=http://localhost:8000
Environment=REDIS_HOST=localhost
Environment=REDIS_PORT=6379
Environment=REDIS_URL=redis://:redis_secret_2024@localhost:6379/0
ExecStartPre=/bin/sleep 30
ExecStart=/opt/conecta-plus/scripts/warmup-service.sh
StandardOutput=append:/opt/conecta-plus/logs/warmup-service.log
StandardError=append:/opt/conecta-plus/logs/warmup-service.log
Restart=no
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
